# systemd/.config/user/battery_notify.service
[Unit]
Description=Battery Low Notification

[Service]
Type=oneshot
ExecStart=%h/.config/i3/battery_notify.sh

